# pkg/xps13/longarmz-xps13.install
post_install() {
	systemctl enable --now tlp.service
	systemctl enable --now bluetooth.service
}
