# pkg/desktop/longarmz-desktop.install
post_install() {
	systemctl enable sddm.service
	systemctl enable --now bluetooth.service
	systemctl --usr enable redshift
	cp /etc/geoclue/longarmz-geoclue.conf /etc/geoclue/geoclue.conf
}
