# pkg/desktop/longarmz-desktop.install
post_install() {
	systemctl --usr enable redshift-gtk
	cp /etc/geoclue/longarmz-geoclue.conf /etc/geoclue/geoclue.conf
}
