# pkg/desktop/longarmz-desktop.install
post_install() {
	systemctl enable --now bluetooth.service
	cp /etc/geoclue/longarmz-geoclue.conf /etc/geoclue/geoclue.conf
}
